# systemd/dropbox.service
[Service]
Type=oneshot
# Create /run/tricorder for the flock lock file
RuntimeDirectory=tricorder
RuntimeDirectoryMode=0755
# Ensure module imports resolve
WorkingDirectory=/apps/tricorder
Environment=PYTHONPATH=/apps/tricorder
ExecStart=/apps/tricorder/bin/run_with_runtime_env.sh --ensure-dirs --dropbox-link /apps/tricorder/dropbox -- /usr/bin/flock -n /run/tricorder/dropbox.lock /apps/tricorder/venv/bin/python3 -m lib.process_dropped_file --scan

[Install]
WantedBy=multi-user.target
