# systemd/voice-recorder.service
[Unit]
Description=Tricorder Live Recorder (WebRTC VAD + RMS)
After=tricorder-audio-restore.service sound.target network-online.target
Wants=tricorder-audio-restore.service network-online.target
PartOf=tricorder.target

[Service]
Type=simple
Environment=PYTHONUNBUFFERED=1
WorkingDirectory=/apps/tricorder
ExecStart=/apps/tricorder/bin/run_with_runtime_env.sh --ensure-dirs -- /apps/tricorder/venv/bin/python -m lib.live_stream_daemon
SyslogIdentifier=live-recorder
Restart=always
RestartSec=3
Nice=5
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
KillMode=process
KillSignal=SIGINT
SendSIGKILL=no
TimeoutStopSec=infinity

[Install]
WantedBy=multi-user.target
