# systemd/web-streamer.service
[Unit]
Description=Tricorder HLS Web Streamer
After=network-online.target voice-recorder.service
Wants=network-online.target
PartOf=tricorder.target

[Service]
Type=simple
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/run/tricorder/runtime.env
ExecStartPre=/apps/tricorder/venv/bin/python -m lib.unit_runtime --env-file /run/tricorder/runtime.env
WorkingDirectory=/apps/tricorder
ExecStart=/apps/tricorder/venv/bin/python -m lib.web_streamer
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
